# etc-systemd-system/second-gunicorn.socket
[Unit]
Description=second gunicorn socket

[Socket]
ListenStream=/run/second-gunicorn.sock

[Install]
WantedBy=sockets.target